# apache.conf
[program:apache]
command=/usr/sbin/httpd -c "User %(ENV_NAGIOS_USER)s" -DFOREGROUND
redirect_stderr=true

