# nagios.conf
[program:nagios]
command=%(ENV_NAGIOS_HOME)s/bin/nagios %(ENV_NAGIOS_HOME)s/etc/nagios.cfg

